# spotweb/rootfs/etc/periodic/hourly/spotweb
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the hourly task to update spots
# ==============================================================================

exec php /app/retrieve.php
