# nzbget/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs the Nginx daemon
# ==============================================================================

# Wait for app to become available
bashio::net.wait_for 46836 localhost 300

# Fix not working web UI with Ingress
# Download the combined js/css files from the nzbget webservice and store locally
wget -q -O "${APP_DIR}/combined.js" "http://localhost:46836/combined.js?lib/jquery.min.js+lib/bootstrap.min.js+lib/raphael.min.js+lib/elycharts.min.js+fasttable.js+index.js+util.js+downloads.js+edit.js+status.js+messages.js+history.js+upload.js+config.js+feed.js"
wget -q -O "${APP_DIR}/combined.css" "http://localhost:46836/combined.css?lib/bootstrap.css+style.css"

bashio::log.info "Starting NGinx..."
exec nginx
